# alloy/rootfs/etc/s6-overlay/s6-rc.d/init-alloy/run
#!/usr/bin/with-contenv bash
# Render /etc/alloy/config.alloy from the add-on options before Alloy starts.
set -e

CONFIG_DIR="/etc/alloy"
CONFIG_FILE="${CONFIG_DIR}/config.alloy"
OPTIONS_FILE="/data/options.json"

mkdir -p "${CONFIG_DIR}" /data/alloy

# --- Read add-on options ---
LOKI_URL=$(jq -r '.loki_url' "${OPTIONS_FILE}")
HOST_LABEL=$(jq -r '.host_label // "homeassistant"' "${OPTIONS_FILE}")
LOG_LEVEL=$(jq -r '.log_level // "info"' "${OPTIONS_FILE}")
ADDITIONAL_CONFIG=$(jq -r '.additional_config // ""' "${OPTIONS_FILE}")

# --- Build relabel drop rules for excluded syslog identifiers ---
# Each entry in exclude_syslog_identifiers becomes a relabel rule that drops
# every journal entry whose syslog identifier matches (anchored regex). Use the
# exact identifier seen in VictoriaLogs, e.g. an add-on slug like
# 'addon_03cabcc9_ring_mqtt'. Dropped lines are still visible in the add-on's
# own log viewer; they just aren't shipped to VictoriaLogs.
DROP_RULES=""
EXCLUDED_IDS=""
while IFS= read -r ident; do
    [ -z "${ident}" ] && continue
    EXCLUDED_IDS="${EXCLUDED_IDS} ${ident}"
    DROP_RULES="${DROP_RULES}
  rule {
    source_labels = [\"__journal_syslog_identifier\"]
    regex         = \"${ident}\"
    action        = \"drop\"
  }"
done < <(jq -r '.exclude_syslog_identifiers[]? // empty' "${OPTIONS_FILE}")

# --- Detect journal path (persistent vs volatile) ---
JOURNAL_PATH="/var/log/journal"
if [ ! -d "${JOURNAL_PATH}" ] || [ -z "$(ls -A "${JOURNAL_PATH}" 2>/dev/null)" ]; then
    JOURNAL_PATH="/run/log/journal"
fi

echo "-----------------------------------------------------------"
echo " Grafana Alloy for Home Assistant (homelab)"
echo " Loki URL:     ${LOKI_URL}"
echo " host label:   ${HOST_LABEL}"
echo " Journal path: ${JOURNAL_PATH}"
echo " Log level:    ${LOG_LEVEL}"
echo " Excluded ids:${EXCLUDED_IDS:- (none)}"
echo "-----------------------------------------------------------"

# The heredoc is unquoted so ${VAR} expands; literal Alloy '$' is written as \$
# and the regex backslash '\s' is written as \\\\s (bash -> \\, Alloy -> \s).
cat > "${CONFIG_FILE}" <<ALLOYCONFIG
// Grafana Alloy config — auto-generated by the HAOS add-on.
// Do not edit directly; configure via the add-on Configuration tab.

logging {
  level = "${LOG_LEVEL}"
}

// --- Read the host systemd journal (mounted via journald: true) ---
loki.source.journal "journal" {
  path          = "${JOURNAL_PATH}"
  max_age       = "12h"
  forward_to    = [loki.process.journal.receiver]
  relabel_rules = loki.relabel.journal.rules
  labels        = {
    host = "${HOST_LABEL}",
    job  = "systemd-journal",
  }
}

// --- Promote useful journal fields to labels ---
// 'host' + 'job' match the homelab cr_alloy convention so HAOS slots into the
// same VictoriaLogs queries/dashboards as every other host.
loki.relabel "journal" {
  forward_to = []

  rule {
    source_labels = ["__journal__systemd_unit"]
    target_label  = "unit"
  }
  rule {
    source_labels = ["__journal__hostname"]
    target_label  = "hostname"
  }
  rule {
    source_labels = ["__journal_syslog_identifier"]
    target_label  = "syslog_identifier"
  }
  rule {
    source_labels = ["__journal__transport"]
    target_label  = "transport"
  }
  rule {
    source_labels = ["__journal_container_name"]
    target_label  = "container_name"
  }
  rule {
    source_labels = ["__journal_priority_keyword"]
    target_label  = "level"
  }${DROP_RULES}
}

// --- Drop empty lines, then ship ---
loki.process "journal" {
  stage.drop {
    expression = "^\\\\s*\$"
  }

  forward_to = [loki.write.victorialogs.receiver]
}

// --- Push to VictoriaLogs (Loki ingestion protocol) ---
loki.write "victorialogs" {
  endpoint {
    url = "${LOKI_URL}"
  }
}
ALLOYCONFIG

# --- Append optional raw Alloy config from the add-on options ---
if [ -n "${ADDITIONAL_CONFIG}" ] && [ "${ADDITIONAL_CONFIG}" != "null" ]; then
    {
        echo ""
        echo "// --- Additional user config (from add-on options) ---"
        echo "${ADDITIONAL_CONFIG}"
    } >> "${CONFIG_FILE}"
fi

echo "Alloy config generated at ${CONFIG_FILE}"
