# alloy/rootfs/etc/s6-overlay/s6-rc.d/alloy/run
#!/usr/bin/with-contenv bash
# Long-running Alloy service. Reads the config rendered by init-alloy.
exec /usr/bin/alloy run \
    --server.http.listen-addr=0.0.0.0:12345 \
    --storage.path=/data/alloy \
    /etc/alloy/config.alloy
